# api.sv.conf
[program:sync-engine]
command=/src/docker.start.api.sh
redirect_stderr=true