# api.sv.conf
[program:sync-engine-api]
user=syncengine2
command=/src/docker.start.api.sh
redirect_stderr=true