# api.sv.conf
[program:sync-engine-api]
command=su - syncengine -c /src/docker.start.api.sh
redirect_stderr=true