# engine.sv.conf
[program:sync-engine]
command=su - syncengine -c /src/docker.start.engine.sh
redirect_stderr=true