# engine.sv.conf
[program:sync-engine]
command=/src/docker.start.engine.sh
redirect_stderr=true