# init.sv.conf
[program:sync-engine-init]
command=/src/docker.sync-engine-init.sh
redirect_stderr=true